# systemd/pi-tor-rules.service
[Unit]
Description=Apply iptables rules for Pi Tor transparent proxy
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
Environment=PROJECT_DIR=/opt/pi-tor
EnvironmentFile=-/etc/default/pi-tor
ExecStart=${PROJECT_DIR}/bin/pi-tor-rules.sh apply
ExecStop=${PROJECT_DIR}/bin/pi-tor-rules.sh flush
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
